# crontab
* 0 * * * tmpclean ~/.tmp
* 0 * * * tmpclean ~/.cache --parent
* 0 * * * tmpclean ~/.ansible_async --days 2 --parent
